# deb/install/prerm
#!/bin/bash

systemctl stop waggle-bk-reverse-tunnel.service || true
systemctl disable waggle-bk-reverse-tunnel.service || true
