# deb/reverse/postinst
#!/bin/bash

case "${1}" in
  configure)
    echo "Enabling & starting waggle-bk-reverse-tunnel service"
    systemctl enable waggle-bk-reverse-tunnel.service
    ;;
esac
